# cron/naver_special_deals.cron
SHELL=/bin/zsh
PATH=/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin

# Crawl Naver Shopping special deals every 30 minutes.
*/30 * * * * /Users/andylee/dev/990won/.venv/bin/python /Users/andylee/dev/990won/scripts/naver_special_deals.py --output-dir /Users/andylee/dev/990won/data/naver_special_deals >> /Users/andylee/dev/990won/logs/naver_special_deals.log 2>&1
