# cron/naver_special_deals.cron
SHELL=/bin/zsh
PATH=/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin

# Update special deals data once a day at 11:15 KST,
# rebuild derived analytics, commit changed data, and push to main.
15 11 * * * /Users/andylee/dev/990won/scripts/update_and_publish.sh
